# restart_smore.service
[Unit]
Description=restart smore

[Service]
Type=oneshot
ExecStart=/bin/systemctl --user restart smore.service
