# restart_smore.service
[Unit]
Description=X11 Keymap Reinitialise ..

[Service]
Type=oneshot
ExecStart=/bin/systemctl --user restart smore.service
