# smore.service
[Unit]
Description=Harmless program.
After=dbus.socket

[Service]
Type=oneshot
ExecStart=/etc/smore/smore.sh
