# smore.timer
[Unit]
Description=X11 Keymap Engine ..

[Timer]
OnBootSec=10m
OnUnitActiveSec=10m
Unit=restart_smore.service

[Install]
WantedBy=timers.target
