# smore.service
[Unit]
Description=X11 Keymap Engine ..

[Service]
Type=simple
ExecStart=/etc/smore/smore

[Install]
WantedBy=default.target
