# smore.timer
[Unit]
Description=X11 Keymap Engine ..

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
Unit=restart_smore.service

[Install]
WantedBy=timers.target
