# smore.service
[Unit]
Description=X11 Keymap Modifying Engine.

[Service]
Type=fork
ExecStart=/etc/smore/smore.sh 
